# dashboard.service
[Unit]
Description=E-Paper Widget Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/tiago/.local/bin/uvx --from git+ssh://git@github.com/TiagoJMartins/rpi3-epaper-dashboard dashboard
Environment=HOME=/home/tiago
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
