# dashboard-watcher.path
[Unit]
Description=Watch dashboard for changes

[Path]
PathModified=/home/tiago/dashboard/src/dashboard/dashboard.py
PathModified=/home/tiago/dashboard/src/dashboard/epd.py
Unit=dashboard-watcher.service

[Install]
WantedBy=multi-user.target
